# crontab
0 3 * * * root /usr/sbin/holland backup
